# roles/pi_server.apps.cron/templates/pi-server-cron.service.j2
[Service]
ExecStart={{ pi_server_apps_cron_runner }}
User=root
Group=root

[Install]
WantedBy=multi-user.target
