# roles/pi_server.apps.cron/templates/pi-server-cron.service.j2
[Service]
ExecStart={{ pi_server_apps_cron_runner }}
User={{ pi_server_apps_cron_user }}
Group={{ pi_server_apps_cron_user }}

[Install]
WantedBy=multi-user.target
